# ansible/roles/httpserver/files/httpserver.conf
[program:httpserver]
environment=LD_LIBRARY_PATH="/usr/local/lib"
command=/server
autostart=true
stderr_logfile = /var/log/httpserver.err
stdout_logfile = /var/log/httpserver.log